# IaC/4-fis-experiment/ansible/roles/etcd-backup/templates/etcd-backup.service.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=etcd snapshot backup to S3
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/etcd-backup.sh
