# IaC/3-v3/ansible/roles/etcd-backup/templates/etcd-backup.timer.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=etcd backup timer (every {{ etcd_backup_interval }})

[Timer]
OnBootSec=15min
OnUnitActiveSec={{ etcd_backup_interval }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
